# S05/ticktock.service
[Unit]
Description=Tick Tock

[Service]
ExecStart=/usr/sbin/ticktock

[Install]
WantedBy=multi-user.target